# relcom/ms14/crontab
# %MAILSERV START%
33       4      * * *   mailserv        /usr/spool/mailserv/bin/backup
# %MAILSERV END%
